# flicd/rootfs/etc/services.d/flicd/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: flicd
# Run the Flic Daemon
# ==============================================================================


bashio::log.info "Starting flicd..."

HCI_DEV=$(bashio::config 'hci_dev')
FLICD_PORT=$(bashio::addon.port 5551)

exec /usr/bin/flicd -f /data/flic-db -s "0.0.0.0" -p "$FLICD_PORT" "$HCI_DEV"
