# photons-interactor/rootfs/etc/services.d/interactor/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: Photons Interactor
# ==============================================================================
ip_address=$(bashio::network.ipv4_address)
port=$(bashio::addon.port 6100)
hostname=$(bashio::host.hostname)

cd /config/photons || bashio::exit.nok "Could not access Photons Interactor configuration directory"

if ! bashio::fs.file_exists '/config/photons/interactor.yml'; then
  bashio::var.json \
          interface "$ip_address" \
          port "$port" \
          hostname "$hostname" \
          | tempio \
              -template /usr/share/tempio/interactor.yml.template \
              -out /config/photons/interactor.yml
fi

export LIFX_CONFIG="/config/photons/interactor.yml"

if bashio::config.true 'override_config_file'; then
  export INTERACTOR_HOST="$ip_address"
  export INTERACTOR_PORT="$port"
  bashio::log.info "Starting Photons Interactor on $hostname ($ip_address:$port)"
else
  bashio::log.info "Starting Photons Interactor using /config/photons/interactor.yml"
fi


if bashio::config.true 'enable_logging'; then
  bashio::log.info "Photons Interactor verbose logging enabled."
  export PHOTONS_SILENT_BY_DEFAULT=0
else
  bashio::log.info "Photons Interactor verbose logging disabled."
  export PHOTONS_SILENT_BY_DEFAULT=1
fi

exec lifx lan:interactor
